# scripts/br-agent.service
[Unit]
Description=bladerunner guest control agent
After=network.target

[Service]
ExecStart=/usr/local/bin/br-agent
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
