# internal/provision/scripts/bladerunner-vsock-ntp.service
[Unit]
Description=Bladerunner vsock NTP forward (guest UDP 123 -> host vsock)
After=network.target
ConditionPathExists=/usr/bin/socat
ConditionPathExists=/dev/vsock

[Service]
Type=simple
ExecStart=/usr/bin/socat UDP4-RECVFROM:123,bind=127.0.0.1,fork,reuseaddr VSOCK-CONNECT:2:18557
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
